# launchd/com.aqara.adapters.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aqara.adapters</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/.aqara/start-adapters.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AQARA_API_KEY</key>
        <string>sk-f12208bbc17b18dfb36e083c08e07ed258e37522841422f1a5a67a0a426fc910</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/.aqara-adapters.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/.aqara-adapters.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
